# supervisor/prod.conf
[program:webserver]
directory = /opt/app
command = gunicorn project.wsgi --bind=0.0.0.0:80 --workers=20 --timeout=600 --limit-request-line=8188
